# deploy/supervisor.conf
; supervisor config for the MeetingStore Ops Horizon worker.
;
; Place at:  /etc/supervisord.d/meetingstore-horizon.ini
;            (RHEL/AlmaLinux — Supervisor's include dir)
;            OR /etc/supervisor/conf.d/meetingstore-horizon.conf
;            (Debian/Ubuntu)
;
; Activate:
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl status meetingstore-horizon       (expect RUNNING)
;
; Tail logs:
;   tail -f /var/log/meetingstore-horizon.log

[program:meetingstore-horizon]
process_name=%(program_name)s
command=php /home/stcav/ms.21stcav.com/artisan horizon
autostart=true
autorestart=true
user=stcav
redirect_stderr=true
stdout_logfile=/var/log/meetingstore-horizon.log
stopwaitsecs=3600
; horizon:terminate sends SIGTERM; the worker drains in-flight jobs before
; exiting, then supervisor respawns it. stopwaitsecs=3600 = 1h grace before
; supervisor escalates to SIGKILL — well above the 600s job timeout ceiling.

; If you ever need to scale: numprocs=2 etc — but Horizon manages its own
; worker pool internally based on config/horizon.php, so a single supervisor
; process is the right shape.
